# modules/scripts/image/.config/systemd/user/lightdark.timer
[Unit]
Description=Automatically set light or dark mode

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
