# modules/scripts/image/.config/systemd/user/lightdark.service
[Unit]
Description=Automatically set light or dark mode

[Service]
Type=oneshot
ExecStart=%h/bin/lightdark
